# files/HOME/.config/autostart/discord.desktop
[Desktop Entry]
Name=Discord
Exec=/usr/bin/discord
Terminal=false
Type=Application
