# dot_config/autostart/signal.desktop
[Desktop Entry]
Name=Signal
Exec=/usr/bin/signal-desktop
Terminal=false
Type=Application
