# dot_config/autostart/thunderbird.desktop
[Desktop Entry]
Name=Thunderbird
Exec=/usr/bin/thunderbird
Terminal=false
Type=Application
