# files/autostart/xbindkeys.desktop
[Desktop Entry]
Name=Xbindkeys
Exec=/usr/bin/xbindkeys
Terminal=false
Type=Application
